# src/systemd/location-update.timer.in
[Unit]
Description=Update the location database once a day

[Timer]
OnCalendar=daily
RandomizedDelaySec=24h

[Install]
WantedBy=timers.target
