# src/systemd/location-update.service.in
[Unit]
Description=Automatic Location Database Updater
Documentation=man:location(1) https://man-pages.ipfire.org/libloc/location.html
Requires=network.target

[Service]
Type=oneshot
ExecStart=@bindir@/location update
