# src/systemd/location-update.service.in
[Unit]
Description=Automatic Location Database Updater
Requires=network.target

[Service]
Type=oneshot
ExecStart=@bindir@/location update
